# resources/container/elophant-db/etc/cron.d/elophant-backup
0 * * * * root /usr/bin/db-maintenance/backup.sh
